# systemd/tarsnap-archive-prune.service
[Unit]
Description=Tarsnap Prune Daily

[Service]
Type=simple
ExecStart=/usr/local/bin/tarsnap-prune.sh
