# systemd/tarsnap-archive-monthly.service
[Unit]
Description=Tarsnap Archive Monthly

[Service]
Type=simple
ExecStart=/usr/local/bin/tarsnap-archive.sh monthly
